# 30/ansible/roles/deploy_app/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
ListenStream={{ gunicorn_bind_address }}:{{ gunicorn_bind_port }}
[Install]
WantedBy=sockets.target
